# config/nodejs.service
[Unit]


[Service]
ExecStart=/opt/nodejs.sh

[Install]
WantedBy=default.target
